# installer/src/raspi_pipgio/work/etc/systemd/system/webapp-weather-finder.service
[Unit]
Description=Flask webapp WeatherFinder service
After=network-online.target

[Service]
Type=idle
ExecStart=/home/pi/webapp/start.sh prod >/dev/null
User=pi

[Install]
WantedBy=multi-user.target

