# installer/src/raspi_pipgio/work/crontab/pi
10 5,10,12,14,16 * * * /home/pi/bin/get_weather_forecast.sh
30 23 * * * /home/pi/bin/remove_applogs.sh
