# installer/src/raspi_disp_oled/work/etc/systemd/system/button_start.service
[Unit]
Description=Button start Specific application subprocess.
After=newtwork-online.target

[Service]
Type=simple
ExecStart=/home/pi/bin/buttons_startup.sh
User=pi

[Install]
WantedBy=multi-user.target

